# launchd/com.example.parkingspottracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Auto-start the parking monitor at login and keep it running.

  Setup:
    1. Edit the two absolute paths below to match your checkout.
    2. cp this file to ~/Library/LaunchAgents/
    3. launchctl load -w ~/Library/LaunchAgents/com.example.parkingspottracker.plist

  caffeinate -i keeps the Mac from idle-sleeping while the monitor runs.
  Logs (incl. detections) go to monitor.log in the project directory.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.parkingspottracker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-i</string>
        <string>/ABSOLUTE/PATH/TO/ParkingSpotTracker/.venv/bin/python</string>
        <string>monitor.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/ABSOLUTE/PATH/TO/ParkingSpotTracker</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>15</integer>

    <key>StandardOutPath</key>
    <string>/ABSOLUTE/PATH/TO/ParkingSpotTracker/monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/ABSOLUTE/PATH/TO/ParkingSpotTracker/monitor.log</string>
</dict>
</plist>
